# docker/php/cron/crontab
* * * * * /usr/local/bin/php /var/www/api/artisan schedule:run --no-ansi > /dev/null 2>> /var/log/cron.log
